# final.conf
[program:final]
command=/bin/bash /root/final2.sh
autostart=true
autorestart=false
stdout_logfile=/var/log/final.out.log
stderr_logfile=/var/log/final.err.log
