# watchdog.conf
[program:watchdog]
command=/bin/bash /root/watchdog2.sh
autostart=true
autorestart=false
stderr_logfile=/var/log/watchdog.err.log
stdout_logfile=/var/log/watchdog.out.log
